# indian-market-agent/deploy/lightsail/market-desk.service
[Unit]
Description=India Market Desk Gunicorn Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/srv/indian-news-bot/indian-market-news-bot
EnvironmentFile=/srv/indian-news-bot/indian-market-news-bot/deploy/lightsail/market-desk.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/srv/indian-news-bot/indian-market-news-bot/.venv/bin/gunicorn backend.app:app --bind 127.0.0.1:8000 --workers 1 --threads 8 --timeout 120
Restart=always
RestartSec=5
TimeoutStopSec=30
SyslogIdentifier=market-desk
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
