# indian-market-agent/deploy/lightsail/market-desk-worker.service
[Unit]
Description=India Market Desk Background Worker
After=network-online.target market-desk.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/srv/indian-news-bot/indian-market-news-bot
EnvironmentFile=/srv/indian-news-bot/indian-market-news-bot/deploy/lightsail/market-desk.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/srv/indian-news-bot/indian-market-news-bot/.venv/bin/python -m backend.worker
Restart=always
RestartSec=5
TimeoutStopSec=30
SyslogIdentifier=market-desk-worker
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
